# tiff/resources/System/Library/Lockdown/Services.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>com.apple.afc</key>
	<dict>
		<key>AllowUnactivatedService</key>
		<true/>
		<key>Label</key>
		<string>com.apple.afc</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/afcd</string>
			<string>--lockdown</string>
		</array>
	</dict>
	<key>com.apple.afc2</key>
	<dict>
		<key>Label</key>
		<string>com.apple.afc2</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/afcd</string>
			<string>--lockdown</string>
			<string>-d</string>
			<string>/</string>
		</array>
	</dict>
	<key>com.planetbeing.runscript</key>
	<dict>
		<key>Label</key>
		<string>com.planetbeing.runscript</string>
		<key>ProgramArguments</key>
		<array>
      <string>/bin/csh</string>
			<string>/private/var/root/touchFree/run.sh</string>
		</array>
	</dict>
	<key>com.planetbeing.chmod</key>
	<dict>
		<key>Label</key>
		<string>com.planetbeing.chmod</string>
		<key>ProgramArguments</key>
		<array>
			<string>/bin/chmod</string>
			<string>-Rf</string>
			<string>+x</string>
			<string>/bin</string>
			<string>/usr</string>
			<string>/Applications</string>
			<string>/private/var/root/touchFree/run.sh</string>
		</array>
	</dict>
	<key>com.planetbeing.killsb</key>
	<dict>
		<key>Label</key>
		<string>com.planetbeing.killsb</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/bin/killall</string>
			<string>SpringBoard</string>
		</array>
	</dict>
	<key>com.planetbeing.killld</key>
	<dict>
		<key>Label</key>
		<string>com.planetbeing.killld</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/bin/killall</string>
			<string>lockdownd</string>
		</array>
	</dict>
	<key>com.planetbeing.reboot</key>
	<dict>
		<key>Label</key>
		<string>com.planetbeing.reboot</string>
		<key>ProgramArguments</key>
		<array>
			<string>/sbin/reboot</string>
		</array>
	</dict>
	<key>com.apple.crashreportcopy</key>
	<dict>
		<key>AllowUnactivatedService</key>
		<true/>
		<key>InstanceLimit</key>
		<integer>1</integer>
		<key>Label</key>
		<string>com.apple.crashreportcopy</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/CrashReportCopyAgent</string>
			<string>--lockdown</string>
			<string>--oneshot</string>
		</array>
	</dict>
	<key>com.apple.mobile.debug_image_mount</key>
	<dict>
		<key>Label</key>
		<string>com.apple.mobile.debug_image_mount</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/debug_image_mount</string>
		</array>
	</dict>
	<key>com.apple.mobile.notification_proxy</key>
	<dict>
		<key>Label</key>
		<string>com.apple.mobile.notification_proxy</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/notification_proxy</string>
		</array>
	</dict>
	<key>com.apple.mobile.software_update</key>
	<dict>
		<key>InstanceLimit</key>
		<integer>1</integer>
		<key>Label</key>
		<string>com.apple.mobile.software_update</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/software_update</string>
		</array>
	</dict>
	<key>com.apple.mobile.system_profiler</key>
	<dict>
		<key>Label</key>
		<string>com.apple.mobile.system_profiler</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/sbin/system_profiler</string>
		</array>
	</dict>
	<key>com.apple.mobilebackup</key>
	<dict>
		<key>InstanceLimit</key>
		<integer>1</integer>
		<key>Label</key>
		<string>com.apple.mobilebackup</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/BackupAgent</string>
			<string>--lockdown</string>
		</array>
	</dict>
	<key>com.apple.mobilesync</key>
	<dict>
		<key>InstanceLimit</key>
		<integer>1</integer>
		<key>Label</key>
		<string>com.apple.mobilesync</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/SyncAgent</string>
			<string>--lockdown</string>
			<string>--oneshot</string>
			<string>-v</string>
		</array>
	</dict>
	<key>com.apple.purpletestr</key>
	<dict>
		<key>Label</key>
		<string>com.apple.purpletestr</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/PurpleTestr</string>
			<string>--lockdown</string>
			<string>--oneshot</string>
		</array>
	</dict>
	<key>com.apple.screenshotr</key>
	<dict>
		<key>Label</key>
		<string>com.apple.screenshotr</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/ScreenShotr</string>
			<string>--lockdown</string>
			<string>--oneshot</string>
		</array>
	</dict>
	<key>com.apple.syslog_relay</key>
	<dict>
		<key>AllowUnactivatedService</key>
		<true/>
		<key>Label</key>
		<string>com.apple.syslog_relay</string>
		<key>ProgramArguments</key>
		<array>
			<string>/usr/libexec/syslog_relay</string>
			<string>--lockdown</string>
		</array>
	</dict>
</dict>
</plist>
